# launchd/com.deng.chrome-bridge-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.deng.chrome-bridge-mcp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Google Chrome.app/Contents/MacOS/Google Chrome</string>
    <string>--remote-debugging-address=127.0.0.1</string>
    <string>--remote-debugging-port=9222</string>
    <string>--user-data-dir=/Users/deng/runtime/.chrome-bridge-mcp/ChromeProfile</string>
    <string>about:blank</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/deng/runtime/.chrome-bridge-mcp/logs/chrome.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/deng/runtime/.chrome-bridge-mcp/logs/chrome.err.log</string>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
